# launchagents/com.zxlab.meta-dashboard-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zxlab.meta-dashboard-server</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>{HOME}/.operacao-ia/scripts/meta/start_dashboard.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>{HOME}/.operacao-ia/logs/meta-dashboard.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>{HOME}/.operacao-ia/logs/meta-dashboard.stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>PORT</key>
    <string>8888</string>
  </dict>
</dict>
</plist>
